# scripts/mejona-blog-service.service
[Unit]
Description=Blog CRM Management Microservice
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/mejona/blog-crm-service
ExecStart=/opt/mejona/blog-crm-service/blog-crm-service
EnvironmentFile=/opt/mejona/blog-crm-service/.env
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target